# email-to-matrix/systemd-files/email-to-matrix.socket
[Socket]
ListenFIFO=%t/email-to-matrix.stdin
Service=email-to-matrix.service
